# systemd/system/tailscale-link-up.path
[Path]
PathChanged=/var/run/tailscale/tailscaled.sock
Unit=tailscale-link-up.service

[Install]
WantedBy=multi-user.target
